# files/etc/profile.d/jython.sh
export PATH=$PATH:/opt/jython/bin
export JYTHON_HOME=/opt/jython
#export JYTHON_OPTS=
